# systemd/chatbot-control-panel.service
[Unit]
Description=Chat Bot Control Panel
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/rasp/chat_bot
EnvironmentFile=/home/rasp/chat_bot/.env
ExecStart=/usr/bin/python3 /home/rasp/chat_bot/control_panel.py
Restart=on-failure
RestartSec=5
User=rasp

[Install]
WantedBy=multi-user.target
